# api/build/docker/queue/resqueue.conf
[program:queue]
command=/usr/local/bin/php /var/slim/bin/queue-worker.php
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true